# deploy/systemd/https-proxy.service
[Unit]
Description=HTTPS Proxy Server
After=network.target

[Service]
Type=simple
User=https-proxy
Group=https-proxy
WorkingDirectory=/opt/https-proxy
ExecStart=/opt/https-proxy/https-proxy --config /etc/https-proxy/config.json
Restart=on-failure
RestartSec=5

# 增加文件描述符限制，提高并发连接能力
LimitNOFILE=1048576

# 增加进程数限制
LimitNPROC=65536

# CPU和内存优化
CPUSchedulingPolicy=fifo
CPUSchedulingPriority=99
IOSchedulingClass=realtime
IOSchedulingPriority=0

# 为代理服务分配更多内存
MemoryLow=512M
MemoryMin=256M

# Security enhancements
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths=/opt/https-proxy/stats /opt/https-proxy/data

[Install]
WantedBy=multi-user.target 